# research/m_z/systemd/road_tutorial/echo@.socket
# <echo@.service>
[Unit]
Description=Echo server service

[Service]
ExecStart=/path/to/socketthing.py
StandardInput=socket
