# a_misc_research/m_z/systemd/logging/mytest1.service
[Unit]
Description=My test service

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/etc/mytests/bash_to_journald.sh
